# luge-relay.service
[Unit]
Description=Luge Relay Timing System
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=relay
WorkingDirectory=/home/relay/LugeRelay
ExecStart=/usr/bin/python3 /home/relay/LugeRelay/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Audio system recovery
ExecStartPost=/bin/sleep 2
ExecStartPost=/home/relay/LugeRelay/audio_monitor.sh check
ExecStartPost=/bin/bash -c 'if [ $? -ne 0 ]; then /home/relay/LugeRelay/audio_monitor.sh reinit; fi'

[Install]
WantedBy=multi-user.target
